# deploy/systemd-run-binary.service
[Unit]
Description=Wartracker Discord Bot (native binary)
After=network.target

[Service]
User=wartracker
Group=wartracker
WorkingDirectory=/opt/wartracker
Environment=BOT_TOKEN=your_token_here
ExecStart=/usr/local/bin/wartracker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
